# scripts/com.scubamount.nim-ranker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.scubamount.nim-ranker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/andrew/code/scubamount-nim-proxy/.venv/bin/python</string>
    <string>-m</string>
    <string>auto_ranker.probe</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/andrew/code/scubamount-nim-proxy</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/nim_ranker.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/nim_ranker.log</string>
</dict>
</plist>
