# scripts/com.scubamount.nim-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.scubamount.nim-proxy</string>
  <key>ProgramArguments</key>
  <array>
    <string>__VENV_PY__</string>
    <string>-u</string>
    <string>-m</string>
    <string>uvicorn</string>
    <string>nim_proxy:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>__PORT__</string>
    <string>--no-access-log</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__PROXY_DIR__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__LOG__</string>
  <key>StandardErrorPath</key>
  <string>__LOG__</string>
</dict>
</plist>
